# deploy/fundmgr-global-web.service
[Unit]
Description=AI Fund Manager — Global Simulation Web Dashboard
After=network.target

[Service]
Type=simple
User=alexander
WorkingDirectory=/home/alexander/Documents/ai-fund-manager
EnvironmentFile=/home/alexander/Documents/ai-fund-manager/.env
Environment=FUND_CONFIG=/home/alexander/Documents/ai-fund-manager/config/config_global.yaml
ExecStart=/home/alexander/Documents/ai-fund-manager/.venv/bin/uvicorn fundmgr.web.app:app --host 0.0.0.0 --port 8011
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
