# deploy/fundmgr-bot.service
[Unit]
Description=AI Fund Manager — Telegram Bot
After=network.target

[Service]
Type=simple
User=alexander
WorkingDirectory=/home/alexander/Documents/ai-fund-manager
EnvironmentFile=/home/alexander/Documents/ai-fund-manager/.env
ExecStart=/home/alexander/Documents/ai-fund-manager/.venv/bin/python -m fundmgr.notify.telegram_bot
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
